# .config/systemd/user/duckdns-update@.timer
[Unit]
Description=Timer for DuckDNS domain updater for %i

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
