# .config/systemd/user/gpg-agent.socket.d/override.conf
[Socket]
ListenStream=
#ListenStream=%t/gnupg/d.7w4n7zasdwuoxdin9nqmqkcu/S.gpg-agent
ListenStream=/run/user/10000/gnupg/d.7w4n7zasdwuoxdin9nqmqkcu/S.gpg-agent
