# .config/containers/systemd/firecrawl-redis.container
[Container]
ContainerName=firecrawl-redis
Image=docker.io/library/redis:alpine
Pod=firecrawl.pod
